# scripts/launchd/com.robssolar.frontend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.robssolar.frontend</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/robertcashman/robs-solar/scripts/frontend-service.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>ROBS_SOLAR_ROOT</key>
        <string>/Users/robertcashman/robs-solar</string>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/robertcashman/robs-solar/frontend</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>StandardOutPath</key>
    <string>/Users/robertcashman/robs-solar/.frontend-service.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/robertcashman/robs-solar/.frontend-service.log</string>
</dict>
</plist>
